# c0ntempt/live-build-config/config/includes.chroot/etc/xdg/autostart/vidalia-autostart.desktop
[Desktop Entry]
Type=Application
Name=Vidalia
Exec=vidalia
X-GNOME-Autostart-enabled=true
X-GNOME-AutoRestart=true
OnlyShowIn=GNOME;Unity;

